# deploy/supervisor/econtract.conf
; ============================================================
; E-Contract Service – Supervisor configuration
; Install: sudo cp econtract.conf /etc/supervisor/conf.d/
;          sudo supervisorctl reread && sudo supervisorctl update
; ============================================================

[program:samurai-hub-econtract]
command=/var/www/client-management-app/venv/bin/python -m e_contract_service.app
directory=/var/www/client-management-app
user=www-data

; Inherit the same environment as the main app. Override as needed.
environment=
    FLASK_ENV="production",
    SECRET_KEY="%(ENV_SECRET_KEY)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    SESSION_COOKIE_NAME="session",
    SESSION_COOKIE_SECURE="true",
    SESSION_COOKIE_SAMESITE="Lax"

autostart=true
autorestart=true
startretries=5
startsecs=3

stdout_logfile=/var/log/supervisor/samurai-hub-econtract.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/samurai-hub-econtract-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
