# deploy/supervisor/gps_tracker.conf
; ============================================================
; GPS Tracker Service – Supervisor configuration
; 車載GPS機器(Teltonika FMM880等)を受け付けるTCPサーバー
;
; Install: sudo cp gps_tracker.conf /etc/supervisor/conf.d/
;          sudo supervisorctl reread && sudo supervisorctl update
;
; 注意: GPS_TCP_PORT(既定5027)をエックスサーバーVPSのパケット
;       フィルターおよびOSファイアウォールで開放すること。
; ============================================================

[program:samurai-hub-gps-tracker]
command=/var/www/client-management-app/venv/bin/python -m gps_tracker
directory=/var/www/client-management-app
user=www-data

; メインアプリと同じDBに接続する。ポート等は必要に応じて変更する。
environment=
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    GPS_TCP_HOST="0.0.0.0",
    GPS_TCP_PORT="5027",
    GPS_TZ_OFFSET_HOURS="9",
    GPS_SOCKET_TIMEOUT="1800"

autostart=true
autorestart=true
startretries=5
startsecs=3

stdout_logfile=/var/log/supervisor/samurai-hub-gps-tracker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/samurai-hub-gps-tracker-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
